# docs/com.open-buddy.voice-helper.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.open-buddy.voice-helper</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/open</string>
        <string>-W</string>
        <string>-n</string>
        <string>-a</string>
        <string>/Users/Shared/.open-buddy/VoiceHelper.app</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/tmp/open-buddy-voice.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/open-buddy-voice.log</string>
    <key>ThrottleInterval</key>
    <integer>5</integer>
</dict>
</plist>
